# scripts/pihole-health.service
[Unit]
Description=Pi-hole Health Logger
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pihole_health_logger.sh
User=root

[Install]
WantedBy=multi-user.target
